# mongodb.conf
[program:mongodb]
command=/usr/bin/mongod --unixSocketPrefix=/var/run/mongodb --config /etc/mongodb.conf --rest --dbpath /data/db
user=root
autostart=true
